# deploy/systemd/aserras-frontend.service
[Unit]
Description=Aserras Frontend FastAPI service
# Generated/managed by deploy/setup.sh. Update the script when making changes.
Documentation=https://github.com/Aserras/Aserras-Frontend
After=network.target

[Service]
Type=simple
User=aserras
Group=aserras
WorkingDirectory=/opt/aserras-frontend
EnvironmentFile=-/opt/aserras-frontend/.env
Environment="PATH=/opt/aserras-frontend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/opt/aserras-frontend/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8080 --proxy-headers --forwarded-allow-ips='*'
Restart=on-failure
RestartSec=5s
StandardOutput=append:/var/log/aserras/frontend.log
StandardError=append:/var/log/aserras/frontend.err.log

[Install]
WantedBy=multi-user.target
